# examples/config/es6_supervisor.conf
# PUT THIS FILE IN /etc/supervisord.conf

[supervisord]
logfile=/data1/logs/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
minfds=100000

[unix_http_server]
file=/etc/supervisor.sock   ; (the path to the socket file)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///etc/supervisor.sock

[program:es]
command=/usr/local/elasticsearch/bin/elasticsearch
directory=/usr/local/elasticsearch
autostart=true
autorestart=true
startretries=10
stopsignal=INT
stopwaitsecs=600
stderr_logfile=/data1/logs/supervisor_es_error.log
stdout_logfile=/data1/logs/supervisor_es.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
user=ec2-user
environment=JAVA_HOME=/usr/java/default

[program:push_to_es]
command=/home/ec2-user/pypy/bin/pypy activedata_etl/push_to_es.py --settings=resources/settings/staging/push_to_es.json
directory=/home/ec2-user/ActiveData-ETL
autostart=true
process_name=%(process_num)02d
numprocs=1
autorestart=true
stopsignal=INT
stopwaitsecs=30
stderr_logfile=/data1/logs/supervisor_push_to_es.error.log
stdout_logfile=/data1/logs/supervisor_push_to_es.log
user=ec2-user
environment=PYTHONPATH=.:vendor;HOME=/home/ec2-user

